# mcp-server.service
[Unit]
Description=My MCP Server
After=network.target

[Service]
Type=simple
User=MCP_USER
Group=MCP_USER
WorkingDirectory=MCP_PROJECT_DIR
Environment=PATH=MCP_PROJECT_DIR/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=MCP_PROJECT_DIR/.venv/bin/python MCP_PROJECT_DIR/main.py --transport http --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-server

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=MCP_PROJECT_DIR
PrivateTmp=true

[Install]
WantedBy=multi-user.target